# project-spec/meta-user/recipes-apps/autostart/assets/systemd/system/autostart.service
[Unit]
Description=Simple autostart application

[Service]
ExecStart=/usr/bin/main -d
ExecReload=/bin/kill -SIGHUP $MAINPID

[Install]
WantedBy=multi-user.target
